# build_context/wls.conf
[program:xe]
command=/bin/bash -c "java -version"
autorestart=true
