# build_context/xe.conf
[program:xe]
command=/bin/bash -c "/etc/init.d/oracle-xe start"
autorestart=false
